# bin/services/py_data_logger.service
[Install]
RequiredBy=multi-user.target

[Unit]
Description=Run data_logger.py indefinitely and on startup.

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/prod_env
ExecStart=/usr/bin/python3 /home/ubuntu/prod_env/data_logger.py comms=True
